# dot_config/autostart/private_ghnotifier.desktop
[Desktop Entry]
Exec=/home/programmer/.local/bin/ghnotifier
Icon=dialog-scripts
Name=ghnotifier
Path=
Type=Application
X-KDE-AutostartScript=true
