# dot_config/autostart/private_xbindkeys.desktop
[Desktop Entry]
Exec=xbindkeys --poll-rc
Icon=
Name=xbindkeys
Path=
Terminal=False
Type=Application
